# systemd-units/gpg-agent@.service
[Unit]
Description=gpg4win to WSL connector for GPG (S.gpg-agent)
Requires=gpg-agent.socket

[Service]
Type=simple
ExecStart=%h/.local/bin/npiperelay.exe -v -p -ep -ei -s -a 'C:/Users/alexm/AppData/Local/gnupg/S.gpg-agent'
StandardInput=socket
StandardError=journal

[Install]
WantedBy=default.target
